# root/usr/lib/systemd/system/univention-container-mode-firstboot.service.d/override.conf
[Unit]
OnFailure=univention-container-mode-firstboot-on-failure.service

[Service]
ExecStartPre=-/usr/bin/touch /dev/shm/univention-container-mode-firstboot-on-failure.service
ExecStartPost=-/usr/bin/rm --force /dev/shm/univention-container-mode-firstboot-on-failure.service
