# root/usr/lib/systemd/system/sssd.service.d/override.conf
[Unit]
OnFailure=sssd-on-failure.service
